# linux/systemd/pasori-check.service
[Unit]
Description=PaSoRi check service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pasori-check.sh
